# dot-config/sv/conky/run
#!/bin/sh

exec conky -c "${XDG_CONFIG_DIR:-$HOME/.config}/conky/$(hostname)-conky.conf"
