# dot-config/sv/runsvdir-davide/run
#!/bin/sh

export USER="davide"
export HOME="/home/davide"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.local/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
