# dot-config/sv/redshift/run
#!/bin/sh

[ -r ./conf ] && . ./conf
export DISPLAY=:0
export XAUTHORITY="$HOME/.config/Xauthority"
exec redshift ${OPTS:-} 2>&1
# exec redshift -v
