# dot-config/sv/trayer/run
#!/bin/sh

exec trayer --edge bottom --align right --SetDockType true --SetPartialStrut true --expand true --width 2 --transparent true --alpha 0.5 --height 33 --tint 0x2F0B3A
