# dot-config/sv/udiskie/run
#!/bin/sh

exec udiskie
# exec udiskie \
#     --no-automount \
#     --smart-tray \
#     --no-file-manager \
#     --terminal='alacritty -e' \
#     --password-cache 10
