# dot-config/sv/beep/run
#!/bin/sh
snooze -H /1 sh -c 'echo "\a" >/dev/tty5' # beep
