# dot-config/sv/runsvdir-user/finish
#!/bin/sh

sv -w600 force-stop /home/$USER/service/*
sv exit /home/$USER/service/*
