# dot-config/sv/aw-server/run
#!/bin/sh

exec $HOME/Packages/activitywatch-mod/aw-server/aw-server
