# dot-config/sv/aw-server/run
#!/bin/sh
sleep 3m
exec "$HOME/Packages/activitywatch/aw-server/aw-server"
