# dot-config/sv/aw-server/run
#!/bin/sh
sleep 60
exec "$HOME/Packages/activitywatch-mod/aw-server/aw-server"
