# dot-config/sv/syncthing/run
#!/bin/sh
sleep 20m
exec syncthing serve --no-browser
