# dot-config/sv/emacs-daemon/run
#!/bin/sh

exec emacs -u $USER --chdir=$HOME --fg-daemon
