# dot-config/sv/emacs-daemon/run
#!/bin/sh
/bin/emacs -u davide --fg-daemon
