# dot-config/sv/emacs-daemon/run
#!/bin/sh
# sleep 10
exec emacs -u "$USER" --chdir="$HOME" --fg-daemon
