# dot-config/sv/emacs-daemon/run
#!/bin/sh

[ -r ./conf ] && . ./conf
exec emacs -u $USER --chdir=$HOME --fg-daemon ${OPTS:-}
