# sv/emacs/run
#!/bin/sh
exec chpst -u davide emacs -u davide --fg-daemon 2>&1
