# dot-config/sv/mpd/run
#!/bin/sh

sv check "${XDG_CONFIG_HOME:-$HOME/.config}/sv/pipewire" >/dev/null || exit 1
exec mpd --no-daemon
