# dot-config/sv/unclutter/run
#!/bin/sh

export DISPLAY=:0
export XAUTHORITY="$HOME/.config/Xauthority"
exec unclutter
