# dot-config/sv/unclutter/run
#!/bin/sh
exec unclutter --jitter 10 --exclude-root --ignore-scrolling
