# dot-config/sv/aww-afk/run
#!/bin/sh
sleep 3m
sv check "${XDG_CONFIG_HOME:-$HOME/.config}/sv/aw-server" >/dev/null || exit 1
exec "$HOME/Packages/activitywatch/aw-watcher-afk/aw-watcher-afk"
